# deploy/sglang/sglang.service
[Unit]
Description=SGLang Server (L40S)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/changye/Coding/llm/deploy/sglang
ExecStart=/home/changye/Coding/llm/deploy/sglang/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
